# ACTUAL_NEW_Yii2_example/docker/supervisor/services.ini
[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=0
loglevel=info
pidfile=/run/supervisord.pid

[program:memcached]
directory=/usr/local/
command=/usr/bin/memcached -u memcached -v -m "%(ENV_MEMCACHED_MEMORY)s" -p 11211 -c "%(ENV_MEMCACHED_MAX_CONNECTIONS)s" -I "%(ENV_MEMCACHED_MAX_ITEM_SIZE)s"
autostart=true
autorestart=true
priority=2
redirect_stderr=true

[program:redis]
command=redis-server
autostart=true
autorestart=true
priority=2
redirect_stderr=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
redirect_stderr=true

[program:crond]
command=/usr/sbin/crond -f -d 8
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:websocket]
command=/services/websocket.sh
autostart=true
autorestart=true
priority=6
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:parsing-consumer]
command=/services/parsing-consumer.sh
autostart=true
autorestart=true
priority=6
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0